# deploy/openproxy.service
[Unit]
Description=OpenPortal API server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=openportal
Group=openportal
WorkingDirectory=/opt/openportal
EnvironmentFile=/opt/openportal/.env
ExecStart=/opt/openportal/openportal
Restart=always
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/openportal /etc/wireguard

# WireGuard requires NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openportal

[Install]
WantedBy=multi-user.target
